# systemd/updatedb.service
[Unit]
Description=UpdateDB Service

[Service]
Type=oneshot
ExecStart=/usr/bin/updatedb

[Install]
WantedBy=default.target
