# systemd/updatedb.timer
[Unit]
Description=UpdateDB Timer

[Timer]
OnStartupSec=1min
OnUnitActiveSec=5min

[Install]
WantedBy=timers.target
